# pingbot.service
[Unit]
Description=Telegram Moderation Bot
After=network.target

[Service]
WorkingDirectory=/opt/pingbot
ExecStart=/usr/bin/python3 run.py
Restart=always
EnvironmentFile=/opt/pingbot/.env

[Install]
WantedBy=multi-user.target
